# src/systemd-system/serial-getty@ttyGS0.service.d/usbgadget.conf
[Unit]
After=usb-gadget.target
ConditionPathExists=/dev/ttyGS0
