# src/systemd-system/usb-gadget-base.service.in
[Unit]
Description=USB Gadget - Base Configuration
After=sys-kernel-config.mount
ConditionPathIsDirectory=/sys/kernel/config/usb_gadget

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=PREFIX/libexec/usb-gadget-scripts/composite-base start
ExecStop=PREFIX/libexec/usb-gadget-scripts/composite-base stop
